# docker-crond-go/cron.d/my-crond
*/2 * * * *   echo mike
*/5 * * * *   echo test crond